# startup/duckiehunt.service
[Unit]
Description=Duckiehunt

[Service]
ExecStart=/home/thfalgou/duckiehunt/startup/duckiehunt-start.sh

[Install]
WantedBy=multi-user.target